# eth0-only/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

sleep 10
ifdown wan # J'ai eu besoin de faire ça pour être sûr d'éviter l'ip privée sur le wan, depuis la maj 18.02.
sh /etc/vlanprio.sh
ifup wan
exit 0
