# 22-03/etc - wan name interface = eth0/crontabs/root
00 * * * * /etc/checkinternet.sh
*/5 * * * * /etc/checknftrules.sh
*/5 * * * * /etc/checkroutev6.sh
