# 22-03/etc - wan name interface = eth0/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

echo "boot here $(date)" >> /root/boot-history.txt
exit 0
