# 22-03 with TV/etc - wan name interface = eth0/crontabs/root
00 * * * * /etc/checkinternet.sh
*/5 * * * * /etc/checknftrules.sh
